# install/etc/udev/rules.d/10.sonyautobackup.rules
ENV{DEVTYPE}=="usb_device", SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="054c", RUN+="/bin/bash /media/jetson/data/source/autobackup/sync.sh"